# tap-station.service
[Unit]
Description=NFC Tap Station Service
After=network.target

[Service]
Type=simple
User=dc-nfc-1
WorkingDirectory=/home/dc-nfc-1/nfc-tap-logger
Environment="PATH=/home/dc-nfc-1/nfc-tap-logger/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/dc-nfc-1/nfc-tap-logger/venv/bin/python -m tap_station.main --config /home/dc-nfc-1/nfc-tap-logger/config.yaml
Restart=always
RestartSec=10
StandardOutput=append:/home/dc-nfc-1/nfc-tap-logger/logs/tap-station.log
StandardError=append:/home/dc-nfc-1/nfc-tap-logger/logs/tap-station-error.log

# Restart on failure
StartLimitInterval=200
StartLimitBurst=5

# Run with elevated privileges for GPIO/I2C access
# Alternative: Add user to gpio and i2c groups
# SupplementaryGroups=gpio i2c

[Install]
WantedBy=multi-user.target
